# roles/base/files/users/vgscq/udev_rules/99-dockstations.rules
SUBSYSTEM=="power_supply", ENV{POWER_SUPPLY_ONLINE}=="1", RUN+="/home/vgscq/.config/bin/udev/dock.sh"
SUBSYSTEM=="power_supply", ENV{POWER_SUPPLY_ONLINE}=="0", RUN+="/home/vgscq/.config/bin/udev/dock.sh"
